# v2g-liberty/rootfs/etc/s6-overlay/s6-rc.d/init-appdaemon/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: V2G Liberty
# Configures V2G Liberty
# ==============================================================================

bashio::log.info "Configuring V2G Liberty..."

# Create directory where log files will be written
mkdir -p /config/logs

# Remove dev-only files that earlier versions copied into the persistent
# /config. The copy below only adds/overwrites, so without this these stale
# dev files would survive an upgrade. They are no longer in the image.
# Introduced in 0.8.1. Can be removed once no installation upgrades from a
# version below 0.8.1 anymore (HA jumps straight to the latest version, so a
# 0.8.0 -> later install would otherwise skip this cleanup).
rm -rf /config/apps/dev_tools \
    /config/tests \
    /config/appdaemon.devcontainer.yaml

# Copy python files
cp -aR /root/appdaemon/* /config \
    || bashio::exit.nok 'Failed to copy V2G-Liberty app'
bashio::log.info "AppDaemon app files copied to /config"

#Create packages folder if it does not exist yet
if ! bashio::fs.file_exists '/homeassistant/packages/v2g_liberty'; then
    mkdir -p /homeassistant/packages/v2g_liberty
fi
# Copy files to packages folder
cp -aR /root/homeassistant/packages/v2g_liberty/* /homeassistant/packages/v2g_liberty \
    || bashio::exit.nok 'Failed to copy V2G-Liberty package'
bashio::log.info "Home Assistant package files copied"

# Create folders for plugins if needed
if ! bashio::fs.file_exists '/homeassistant/www/v2g_liberty/apexcharts-card'; then
    mkdir -p /homeassistant/www/v2g_liberty/apexcharts-card
fi
if ! bashio::fs.file_exists '/homeassistant/www/v2g_liberty/lovelace-card-mod'; then
    mkdir -p /homeassistant/www/v2g_liberty/lovelace-card-mod
fi

if ! bashio::fs.file_exists '/homeassistant/www/v2g_liberty/v2g-liberty-cards'; then
    mkdir -p /homeassistant/www/v2g_liberty/v2g-liberty-cards
fi

# Copy plugins
cp -aR /root/homeassistant/www/v2g_liberty/apexcharts-card/* /homeassistant/www/v2g_liberty/apexcharts-card \
    || bashio::exit.nok 'Failed to copy apexchart-card plugin'

cp -aR /root/homeassistant/www/v2g_liberty/lovelace-card-mod/* /homeassistant/www/v2g_liberty/lovelace-card-mod \
    || bashio::exit.nok 'Failed to copy lovelace-card-mod plugin'

cp -aR /root/homeassistant/www/v2g_liberty/v2g-liberty-cards/* /homeassistant/www/v2g_liberty/v2g-liberty-cards \
    || bashio::exit.nok 'Failed to copy v2g-liberty-cards plugin'
bashio::log.info "Front-end plugins copied (apexcharts-card, lovelace-card-mod, v2g-liberty-cards)"

cp -a /root/homeassistant/v2g_liberty_example_configuration.yaml /homeassistant/v2g_liberty_example_configuration.yaml \
    || bashio::exit.nok 'Failed to copy v2g_liberty_example_configuration.yaml'

cp -a /root/homeassistant/quasar_load_balancer.json-example /homeassistant/quasar_load_balancer.json-example \
    || bashio::exit.nok 'Failed to copy quasar_load_balancer.json-example'

bashio::log.info "V2G Liberty configuration complete."
